# systemd/.config/systemd/user/swww-shuffle.timer
[Unit]
Description=Run swww-shuffle periodically

[Timer]
OnUnitActiveSec=1800

[Install]
WantedBy=timers.target
